# setup-logcollector/logyou_theme/update-motd.d/01-logyou-text
#!/bin/sh
blue="\e[1;34m"
reset="\e[0m"
cyan="\e[1;36m"
green="\e[1;32m"
red="\e[1;31m"
pink="\e[95m"
frame_color=$green
text_color=$pink
printf "\n"
printf "$frame_color###########################################################\n"
printf "$frame_color##$text_color  ██╗      ██████╗  ██████╗██╗   ██╗ ██████╗ ██╗   ██╗ $frame_color##\n"
printf "$frame_color##$text_color  ██║     ██╔═══██╗██╔════╝╚██╗ ██╔╝██╔═══██╗██║   ██║ $frame_color##\n"
printf "$frame_color##$text_color  ██║     ██║   ██║██║  ███╗╚████╔╝ ██║   ██║██║   ██║ $frame_color##\n"
printf "$frame_color##$text_color  ██║     ██║   ██║██║   ██║ ╚██╔╝  ██║   ██║██║   ██║ $frame_color##\n"
printf "$frame_color##$text_color  ███████╗╚██████╔╝╚██████╔╝  ██║   ╚██████╔╝╚██████╔╝ $frame_color##\n"
printf "$frame_color##$text_color  ╚══════╝ ╚═════╝  ╚═════╝   ╚═╝    ╚═════╝  ╚═════╝  $frame_color##\n"
printf "$frame_color###########################################################\n"
printf "$reset\n"
printf " * Program:        "$cyan"$(logyou)\n$reset"
printf " * Company:        "$cyan"United Information Highway Co., Ltd.\n$reset"
printf " * Tel:            "$cyan"+66 2016 5678\n$reset"
printf " * Email:          "$cyan"cc_support@uih.co.th\n$reset"
printf "$reset\n"
printf " * Try $green'logyou -i'$reset to check program info\n"
printf "\n"
printf " * This machine IP : $(hostname -I)\n"